# crontab
# m h dom mon dow command
* * * * * /home/uniqueweb/unique-web/update.sh
