# systemd/user/uniqueweb.service
[Service]
ExecStart=/var/www/unique/run.sh
Restart=always

[Install]
WantedBy=default.target
